# deploy/djangoapp.conf
[program:djangoapp]
command=/path/to/mywork/myenv/myproj/deploy/run_gunicorn
directory=/path/to/mywork/myenv/myproj
user=www-data
autostart=true
autorestart=true
redirect_stderr=True
